# claude/scheduled/ruby-app-refactor.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.ruby-app-refactor</string>
  <key>ProgramArguments</key>
  <array>
    <string>__HOME__/.local/bin/ruby-app-refactor</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>RUBY_APP_REFACTOR_REPO</key>
    <string>__HOME__/development/irms-api</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>10</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/claude-scheduled/launchd.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/claude-scheduled/launchd.stderr.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
